# files/etc/profile.d/motd.sh
# FOLLOWING ECHO SUPPORTS USE OF ANSI ESCAPE CODES FOR FORMATTING
echo -e "$(cat /etc/motd.d/*)"
